# deploy/systemd/birdcast-uk-radars-refresh.timer
[Unit]
Description=Refresh UK BirdCast radar metadata daily

[Timer]
OnBootSec=3min
OnCalendar=*-*-* 03:05:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
