# deploy/systemd/birdcast-uk-object-store-plan.timer
[Unit]
Description=Plan UK BirdCast Object Store publication hourly

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
