# deploy/systemd/birdcast-uk-era5-request-smoke.timer
[Unit]
Description=Build UK BirdCast ERA5 request artifact daily

[Timer]
OnBootSec=15min
OnCalendar=*-*-* 05:20:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
