# deploy/systemd/birdcast-uk-historical-build.service
[Unit]
Description=Build UK BirdCast historical VPTS reanalysis artifacts
After=network-online.target birdcast-uk-radars-refresh.service
Wants=network-online.target

[Service]
Type=oneshot
User=birdcast
Group=birdcast
EnvironmentFile=/etc/birdcast-uk/birdcast-uk.env
WorkingDirectory=/opt/birdcast-uk/repo
ExecStart=/opt/birdcast-uk/venv/bin/birdcast-uk historical build --source-dir ${BIRDCAST_UK_HISTORICAL_SOURCE} --output-root ${BIRDCAST_UK_ARTIFACT_ROOT} --radars ${BIRDCAST_UK_RADARS_FILE}

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/opt/birdcast-uk
