# deploy/systemd/birdcast-uk-feature-join.timer
[Unit]
Description=Join UK BirdCast observed and ERA5 features hourly

[Timer]
OnBootSec=9min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
