# deploy/systemd/birdcast-uk-vpts-inventory.timer
[Unit]
Description=Build UK BirdCast VPTS Object Store inventory hourly

[Timer]
OnBootSec=4min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
